# installers/linux/lessoncue-update.service
[Unit]
Description=Run a protected LessonCue server operation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/lessoncue-update
PrivateTmp=true
ProtectHome=true
UMask=0022
