# installers/linux/lessoncue-update.path
[Unit]
Description=Watch for a LessonCue update request

[Path]
PathChanged=/var/lib/lessoncue/config/update-request
Unit=lessoncue-update.service

[Install]
WantedBy=multi-user.target
