# installers/linux/lessoncue.service
[Unit]
Description=LessonCue classroom media server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lessoncue
Group=lessoncue
WorkingDirectory=/opt/lessoncue
ExecStart=/opt/lessoncue/LessonCue.Server
Environment=LESSONCUE_DATA_PATH=/var/lib/lessoncue
# Kept as an upgrade fallback for pre-0.9 installers; fresh installs persist port 80 in the data directory.
Environment=LESSONCUE_HTTP_PORT=8080
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/lessoncue

[Install]
WantedBy=multi-user.target
