# installers/linux/lessoncue-update.service
[Unit]
Description=Install the latest verified LessonCue release
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/lessoncue-update
PrivateTmp=true
ProtectHome=true
UMask=0022
